# systemd/.config/systemd/user/newsboat-reload.service
[Unit]
Description=Auto fetch for newsboat
Wants=network-online.target
After=network-online.target

[Service]
ExecStart=/usr/bin/newsboat -x reload
