# git/.config/systemd/user/git-autofetch.service
[Unit]
Description=Auto fetch for git repo
Wants=network.target

[Service]
ExecStart=sh %h/.local/bin/git-autofetch
