# system/etc/udev/rules.d/10-keyboard-hotplug.rules
# IKBC C87
SUBSYSTEM=="usb",ACTION=="bind",ATTRS{idProduct}=="4042",RUN+="/usr/local/bin/caps2esc"
#ATTRS{idVendor}=="04B4" not working???
