# sfeed/.config/systemd/user/sfeed-fetch.service
[Unit]
Description=Auto fetch for sfeed
Wants=network-online.target
After=network-online.target

[Service]
ExecStart=sfeed_update %h/.config/sfeedrc
